# files/etc/uci-defaults/commotion-splash
#!/bin/sh

[ -f "/usr/share/commotion/files/splash.html" ] && { \
mv -f /usr/share/commotion/files/splash.html /etc/nodogsplash/htdocs/splash.html
}

exit 0